# board/uvic-rocketry/srad-camera-rk3566/rootfs_overlay/etc/systemd/system/usb-gadget.service
[Unit]
Description=USB Gadget Serial Setup
After=local-fs.target
Before=getty.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-usb-gadget.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
